# doc/supervisor/ojms.ini
[group:ojms-job]
programs = user-job
priority = 999

[group:ojms-listen]
programs = listen-enterprise,listen-supplier,listen-user,listen-project,listen-order,listen-log,listen-other
priority = 999

[program:user-job]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600 --sleep=3 --tries=1 --memory=1024
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                          ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-enterprise]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-supplier]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-user]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-project]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-order]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-log]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项

[program:listen-other]
process_name = %(program_name)s
command = php artisan queue:work redis --queue=%(program_name)s --timeout=600  --sleep=3 --tries=1
numprocs = 1
autostart = true
autorestart = true
priority = 999                                      ;优先级 从小到大
startretries = 3                                    ;当进程启动失败后，最大尝试启动的次数
stopwaitsecs = 60                                   ;停止的等待时间
user = www                                              ;请部署同学根据实际情况调整该项
directory = /apps/working/ojms-api                     ;请部署同学根据实际情况调整该项
stdout_logfile = /apps/logs/ojms-api/%(program_name)s.log   ;请部署同学根据实际情况调整该项
stderr_logfile = /apps/logs/ojms-api/%(program_name)s.error ;请部署同学根据实际情况调整该项